# docker/supervisord.conf
[supervisord]
nodaemon=true
user=cliff
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisor/supervisord.pid

[program:cliff]
command=uv run --directory /app/backend uvicorn cliff.main:app
    --host %(ENV_CLIFF_APP_HOST)s
    --port %(ENV_CLIFF_APP_PORT)s
    --workers 1
    --log-level info
directory=/app/backend
environment=CLIFF_STATIC_DIR="/app/frontend/dist"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=3
stopwaitsecs=10
